# debian/raspberrypi-net-mods.service
[Unit]
Description=Copy user wpa_supplicant.conf
ConditionPathExists=/boot/wpa_supplicant.conf
Before=dhcpcd.service NetworkManager.service
After=systemd-rfkill.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/raspberrypi-net-mods/wpa_copy

[Install]
WantedBy=multi-user.target
